# doord.service
[Unit]
Description=Door Daemon

[Service]
ExecStart=/usr/bin/doord
Restart=always

[Install]
WantedBy=multi-user.target
