# install/doord-export.service.template
[Unit]
Description=Door Daemon GPIO Setup

[Service]
Type=oneshot
ExecStart=/PATH/TO/bin/doord export
ExecStop=/PATH/TO/bin/doord unexport
RemainAfterExit=yes
